# system_files/kinoite/usr/etc/profile.d/aurora-fastfetch.sh
#!/usr/bin/sh

alias fastfetch="/usr/bin/fastfetch --file /usr/share/ublue-os/sd-logo.txt --logo-type file -c /usr/share/ublue-os/ublue-os.jsonc"